# systemd/fastmcp-bridge.service
[Unit]
Description=Codex Remote Access Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=src404
WorkingDirectory=/home/src404/src/codex-tunnel
EnvironmentFile=/home/src404/src/codex-tunnel/.env
ExecStartPre=/usr/bin/mkdir -p /home/src404/.local/state/codex/codex-tunnel
ExecStart=/usr/bin/node /home/src404/src/codex-tunnel/server.js
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fastmcp-bridge
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
